# crontab
0 * * * * cd /app && /usr/local/bin/python main.py > /dev/stdout
