# crontab
CONFIG_FILE = /bot/config/prod.ini
0 * * * * /usr/local/bin/python /bot/bot.py >> /bot/cron.log 2>&1
