# crontab
0 * * * * cd /app && python main.py >> cron.log 2>&1